# buildroot-external/board/jethome/jethub-j80/rootfs-overlay/etc/udev/rules.d/05-jethub_set_eth_leds.rules
SUBSYSTEM=="net",KERNEL=="eth0",RUN+="/usr/lib/jethome/jethub_set-eth_leds"
