# buildroot-external/rootfs-overlay/usr/lib/systemd/system/hassos-config.service
[Unit]
Description=JHOS Configuration Manager
After=os-bind.target network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/os-config

[Install]
WantedBy=multi-user.target
