# buildroot-external/rootfs-overlay/usr/lib/systemd/system/mnt-data.mount
[Unit]
Description=OS data partition
Wants=os-expand.service
DefaultDependencies=no
After=os-expand.service
Before=umount.target local-fs.target
Conflicts=umount.target
After=systemd-fsck@dev-disk-by\x2dlabel-os\x2ddata.service
Wants=systemd-fsck@dev-disk-by\x2dlabel-os\x2ddata.service systemd-growfs@mnt-data.service

[Mount]
What=/dev/disk/by-label/os-data
Where=/mnt/data
Type=ext4
Options=commit=30

[Install]
WantedBy=local-fs.target
