# buildroot-external/rootfs-overlay/usr/lib/systemd/system/serial-getty@.service.d/hassos.conf
[Unit]
Wants=os-supervisor.service
After=os-supervisor.service
